# .config/systemd/user/wallee.service
[Unit]
Description=Wallpaper slide
After=lightdm.service

[Service]
Type=simple
ExecStart=/usr/bin/wallee

[Install]
WantedBy=multi-user.target
