# .config/systemd/user/wallee.service
[Unit]
Description=Wallpaper slide

[Service]
Type=simple
ExecStart=/usr/bin/wallee

[Install]
WantedBy=xorg.target
